# mqtt-sensors/bme280.conf.example
[program:bme280]
command=/home/pi/VekotinVerstas/mqtt-sensors/env/bin/python /home/pi/VekotinVerstas/mqtt-sensors/bme280.py
