# mqtt-sensors/sds011.conf.example
[program:sds011]
command=/home/pi/VekotinVerstas/mqtt-sensors/env/bin/python /home/pi/VekotinVerstas/mqtt-sensors/sds011.py /dev/ttyUSB0
